# etc/s6/sv/hwclock/up
envfile /etc/s6/config/hwclock.conf
importas -uD "UTC" HARDWARECLOCK HARDWARECLOCK
ifelse { test $HARDWARECLOCK = UTC } { hwclock --systz -u --noadjfile } hwclock --systz -l --noadjfile
