# etc/s6/sv/hwclock/up
envfile /etc/s6/s6.conf
importas -iu HARDWARECLOCK HARDWARECLOCK
foreground { if { test $HARDWARECLOCK = UTC } hwclock --systz -u --noadjfile }
if -x 0 { test $HARDWARECLOCK = localtime } hwclock --systz -l --noadjfile
