# etc/s6/sv/proc/up
mount -no nodev,noexec,nosuid -t proc proc /proc
