# etc/s6/sv/proc/up
if -n { mountpoint -q /proc }
mount -t proc proc /proc
