# etc/s6/sv/kmod-static-nodes/up
foreground { install -d /run/tmpfiles.d }
kmod static-nodes -f tmpfiles -o /run/tmpfiles.d/kmod.conf
