# etc/s6/sv/cgroups/up
mount -t cgroup2 cgroup2 /sys/fs/cgroup
