# etc/s6/sv/tty5/run
#!/usr/bin/execlineb -P

agetty -8 -L tty5 115200
