# etc/s6/current/run-image/service/s6-linux-init-early-getty/run
#!/usr/bin/execlineb -P

/usr/bin/agetty -8 -a root -L tty1 115200
