# etc/s6/sv/wtmp/up
install -g utmp -m0600 -o utmp /dev/null /var/log/wtmp
