# etc/s6/sv/devpts/up
if { mkdir -p /dev/pts }
mount -o noexec,nosuid,gid=5,mode=600 -t devpts devpts /dev/pts
