# etc/s6/sv/devpts/up
if { install -d /dev/pts }
mount -o noexec,nosuid,gid=5,mode=0600 -t devpts devpts /dev/pts
