# etc/s6/sv/udevadm/up
foreground { udevd -d }
foreground { redirfd -w 1 /dev/null udevadm trigger -c add -a dev }
foreground { redirfd -w 1 /dev/null udevadm trigger -t subsystems -c add }
foreground { redirfd -w 1 /dev/null udevadm trigger -t devices -c add }
foreground { redirfd -w 1 /dev/null udevadm settle }
udevadm control -x
