# etc/s6/sv/udevd-log/run
#!/usr/bin/execlineb -P
foreground { install -dg log -o log /var/log/udevd }
s6-setuidgid log exec -c s6-log -d3 -b -- n3 s6000000 T /var/log/udevd
