# etc/s6/sv/shm/up
if { install -d /run/shm }
mount -no noatime,nodev,noexec,nosuid -t tmpfs shm /run/shm
