# etc/s6/sv/shm/up
if { install -d /run/shm }
ln -fns /run/shm /dev
