# etc/s6/sv/shm/up
foreground { install -d /run/shm }
foreground { mount -no noatime,nodev,noexec,nosuid -t tmpfs shm /run/shm }
ln -fns ../run/shm /dev/shm
