# etc/s6/sv/shm/up
if { mkdir -p /run/shm }
ln -fns /run/shm /dev
