# etc/s6/sv/fstab/up
mount -aO no_netdev -t "nosysfs,nodevpts,noproc,nonfs,nonfs4,nosmbfs,nocifs"
