# etc/s6/sv/fstab/up
foreground { mount -a -O no_netdev }
foreground { zfs mount -a }
