# etc/s6/sv/tmp/up
if { install -dm1777 /run/tmp }
if { mount -no noatime,nodev,nosuid -t tmpfs tmp /run/tmp }
install -dm1733 /run/tmp/.ICE-unix /run/tmp/.X11-unix
