# etc/s6/sv/utmp-early/up
if { mkdir -pm 755 /run/utmps }
if { chown -R 20:20 /run/utmps }
s6-setuidgid utmp
redirfd -w 1 /run/utmps/utmp
exit 0
