# etc/s6/sv/utmp-early/up
if { mkdir -p -m 0755 /run/utmps }
if { chown utmp:utmp /run/utmps }
s6-setuidgid utmp
redirfd -w 1 /run/utmps/utmp
exit 0
