# etc/s6/sv/sys/up
if { mount -no noatime,nodev,noexec,nosuid,ro -t sysfs sys /sys }
if { test -d /sys/firmware/efi }
mount -no ro -t efivarfs efivarfs /sys/firmware/efi/efivars
