# etc/s6/sv/sys/up
if -nt { mountpoint -q /sys }
mount -t sysfs sys /sys
