# etc/s6/sv/tty4/run
#!/usr/bin/execlineb -P

agetty -8 -L tty4 115200
