# etc/s6/sv/tty3/run
#!/usr/bin/execlineb -P

agetty -8 -L tty3 115200
