# etc/s6/current/run-image/service/s6-linux-init-shutdownd/run
#!/usr/bin/execlineb -P

s6-linux-init-shutdownd -d 3 -c /etc/s6/current
