# etc/s6/sv/fd/up
if { ln -fs /proc/self/fd /dev }
if { ln -fs fd/0 /dev/stdin }
if { ln -fs fd/1 /dev/stdout }
if { ln -fs fd/2 /dev/stderr }
