# etc/s6/sv/tty1/run
#!/usr/bin/execlineb -P

agetty -8 -L tty1 115200
