# etc/s6/sv/tty1/run
#!/usr/bin/execlineb -P

agetty -8 -a root -L tty1 115200
