# etc/s6/sv/tty6/run
#!/usr/bin/execlineb -P

agetty -8 -L tty6 115200
