# etc/s6/run-image/service/s6-linux-init-early-getty/run
#!/usr/bin/execlineb -P
agetty -l /usr/bin/login tty1
